# daemontools/service/load_average_monitoring/run
#!/bin/sh
exec 2>&1
exec /cdn77_infra/daemontools/scripts/monitor_load_average.sh
