# daemontools/service/load_average_monitoring/run
#!/bin/sh
exec 2>&1
exec /deployment/daemontools/scripts/monitor_load_average.sh
