# daemontools/service/nginx_stub_status_monitoring/run
#!/bin/sh
exec 2>&1
exec /cdn77_infra/daemontools/scripts/monitor_nginx_stub_status.sh
